# spacontrol.service
[Unit]
Description=Spa Control Web Application
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=spacontrol
Group=spacontrol
WorkingDirectory=/opt/spacontrol
Environment=PATH=/opt/spacontrol/venv/bin
ExecStart=/opt/spacontrol/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spacontrol

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/spacontrol

[Install]
WantedBy=multi-user.target
