# spacontrol.service
[Unit]
Description=Spa Control Web Application
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/spacontrol
Environment=PATH=/home/pi/spacontrol/venv/bin
ExecStart=/home/pi/spacontrol/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spacontrol

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/pi/spacontrol

[Install]
WantedBy=multi-user.target
